# firebase-backend.service
[Unit]
Description=Firebase Email Campaign Backend
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/firebase-manager
Environment=PATH=/var/www/firebase-manager/venv/bin:/usr/bin:/usr/local/bin
Environment=PYTHONPATH=/var/www/firebase-manager
ExecStart=/var/www/firebase-manager/venv/bin/python src/utils/firebaseBackend.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target 