# firebase-backend.service
[Unit]
Description=Firebase Campaign Manager Backend
After=network.target redis.service

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/firebase-manager
ExecStart=/var/www/firebase-manager/venv/bin/python3 src/utils/firebaseBackend.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1
StandardOutput=append:/var/log/firebase-backend.log
StandardError=append:/var/log/firebase-backend.log

[Install]
WantedBy=multi-user.target